# boot/services/inicio.service
[Unit]
After=network.target

[Service]
ExecStart=/home/huan/Scripts/boot/services/at_boot.sh

[Install]
WantedBy=multi-user.target
